# relayer/rootfs/etc/services.d/relayer/finish
#!/usr/bin/with-contenv bash

SOCKET_PATH="${SOCKET_PATH:-/tmp/relayer.sock}"

# delete the socket file if it exists
if [ -S "${SOCKET_PATH}" ]; then
    rm -f "${SOCKET_PATH}"
fi