# relayer/rootfs/etc/services.d/relayer/run
#!/usr/bin/with-contenv bash

DAEMON=/usr/bin/centralized-relay

RLY_CMD="exec $DAEMON"

exec s6-setuidgid relayer /bin/bash -c "$RLY_CMD"
